# private_dot_config/systemd/user/pipewire.service.d/wait-for-usb-audio.conf
[Service]
# USB audio module (snd-usb-audio) loads ~3s after sound.target on this box.
# If pipewire starts before then, it hot-plug-grabs a UCM verb on the USB card
# and later wireplumber UCM probes fail with "Device or resource busy", leaving
# Speaker/SPDIF sinks missing. Block until card is present, then settle.
ExecStartPre=/bin/sh -c 'for i in $(seq 1 50); do grep -q "Generic USB Audio" /proc/asound/cards && break; sleep 0.1; done; sleep 2'
