# private_dot_config/systemd/user/wireplumber.service.d/wait-for-sound.conf
[Unit]
After=sound.target
Wants=sound.target

[Service]
# Wait up to 5s for USB audio card to enumerate before starting wireplumber.
# Without this, wireplumber races the USB device init and the headphone sink
# creation aborts ("PipeWire proxy destroyed"), requiring a manual restart.
ExecStartPre=/bin/sh -c 'for i in $(seq 1 50); do grep -q "Generic USB Audio" /proc/asound/cards && break; sleep 0.1; done; sleep 2'
